# roles/node_exporters/tasks/templates/chrony_exporter.service.j2
[Unit]
Description=Chrony Exporter
Requires=chrony.service

[Service]
ExecStart=/usr/local/bin/chrony_exporter

[Install]
WantedBy=multi-user.target